# etc/supervisor/conf.d/ssh.conf
[program:samba]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
