# etc/supervisor/conf.d/xrdp.conf
[program:xrdp]
command=/usr/sbin/xrdp --nodaemon
autostart=true
autorestart=true
